# rc.local
#!/bin/bash
mount /dev/sda1 /home/share
